# launchd_recovered/jp.openclaw.kaikun04_router_worker_v1.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>KAIKUN04_ROUTER_BOT_TOKEN</key>
		<string>8340634697:AAEAFHxJmqSr2jopCdSyPRZuExuLNvOm4wE</string>
		<key>KAIKUN04_ROUTER_CHAT_ID</key>
		<string>512370342</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>jp.openclaw.kaikun04_router_worker_v1</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-lc</string>
		<string>cd "$HOME/AI/openclaw-factory-daemon" || exit 1; source .venv/bin/activate || exit 1; export OCLAW_DB_PATH="$HOME/AI/openclaw-factory/data/openclaw.db"; exec python -u bots/kaikun04_router_worker_v1.py &gt;&gt; logs/kaikun04_router_worker_v1.out 2&gt;&gt; logs/kaikun04_router_worker_v1.err</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/doyopc/AI/openclaw-factory-daemon/logs/kaikun04_router_worker_v1.launchd.err</string>
	<key>StandardOutPath</key>
	<string>/Users/doyopc/AI/openclaw-factory-daemon/logs/kaikun04_router_worker_v1.launchd.out</string>
	<key>WorkingDirectory</key>
	<string>/Users/doyopc/AI/openclaw-factory-daemon</string>
</dict>
</plist>
