# archive/launchagents_disabled_20260313/jp.openclaw.proposal_auto_approve_v1.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
<key>Label</key><string>jp.openclaw.proposal_auto_approve_v1</string>
<key>ProgramArguments</key>
<array>
<string>/Users/doyopc/AI/openclaw-factory-daemon/scripts/run_proposal_auto_approve_v1.sh</string>
</array>
<key>RunAtLoad</key><true/>
<key>KeepAlive</key><true/>
<key>StandardOutPath</key><string>/Users/doyopc/AI/openclaw-factory-daemon/logs/proposal_auto_approve_v1.launchd.out</string>
<key>StandardErrorPath</key><string>/Users/doyopc/AI/openclaw-factory-daemon/logs/proposal_auto_approve_v1.launchd.err</string>
</dict>
</plist>
