# archive/launchagents_20260316/_broken_20260306-015510/jp.openclaw.spec_refiner_v2.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>jp.openclaw.spec_refiner_v2</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>WorkingDirectory</key><string>/Users/doyopc/AI/openclaw-factory-daemon</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>DB_PATH</key><string>/Users/doyopc/AI/openclaw-factory-daemon/data/openclaw_real.db</string>
    <key>OCLAW_DB_PATH</key><string>/Users/doyopc/AI/openclaw-factory-daemon/data/openclaw_real.db</string>
    <key>PYTHONPATH</key><string>/Users/doyopc/AI/openclaw-factory-daemon</string>
  </dict>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd "/Users/doyopc/AI/openclaw-factory-daemon" && source .venv/bin/activate && set -a && (source env/telegram_daemon.env 2>/dev/null || source env/telegram.env 2>/dev/null || true) && set +a && exec python -u bots/spec_refiner_v2.py</string>
  </array>
  <key>StandardOutPath</key><string>/Users/doyopc/AI/openclaw-factory-daemon/logs/spec_refiner_v2.out</string>
  <key>StandardErrorPath</key><string>/Users/doyopc/AI/openclaw-factory-daemon/logs/spec_refiner_v2.err</string>
</dict>
</plist>
