# archive/launchagents_20260316/jp.openclaw.supervisor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>jp.openclaw.supervisor</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/doyopc/AI/openclaw-factory-daemon/scripts/start_openclaw.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>WorkingDirectory</key>
  <string>/Users/doyopc/AI/openclaw-factory-daemon</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/doyopc</string>
    <key>DB_PATH</key>
    <string>/Users/doyopc/AI/openclaw-factory-daemon/data/openclaw_real.db</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/doyopc/AI/openclaw-factory-daemon/logs/openclaw.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/doyopc/AI/openclaw-factory-daemon/logs/openclaw.err</string>
</dict>
</plist>
